# ir-remote-wizard/rootfs/etc/s6-overlay/s6-rc.d/init-ir-wizard/up
#!/command/execlineb

foreground {
    /etc/s6-overlay/scripts/init-ir-wizard.sh
}
