# ir-remote-wizard/rootfs/etc/s6-overlay/s6-rc.d/ir-wizard/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when the service exits
# ==============================================================================

s6-envuidgid root
/run/s6/basedir/bin/halt
